# chroot/var/lib/dpkg/info/bootlogd.conffiles
/etc/init.d/stop-bootlogd-single
/etc/init.d/stop-bootlogd
/etc/init.d/bootlogd
