# chroot/var/lib/dpkg/info/bootlogd.list
/.
/sbin
/sbin/bootlogd
/etc
/etc/init.d
/etc/init.d/stop-bootlogd-single
/etc/init.d/stop-bootlogd
/etc/init.d/bootlogd
/usr
/usr/share
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/bootlogd
/usr/share/doc
/usr/share/doc/bootlogd
/usr/share/doc/bootlogd/changelog.Debian.gz
/usr/share/doc/bootlogd/copyright
/usr/share/doc/bootlogd/changelog.gz
/usr/share/man
/usr/share/man/man8/bootlogd.8.gz
